# root/etc/services.d/unifi-protect-backup/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc unifi-protect-backup -${VERBOSITY}
