# docker_root/etc/services.d/unifi-protect-backup/run
#!/usr/bin/with-contenv bash


export RCLONE_CONFIG=/config/rclone/rclone.conf
exec \
    s6-setuidgid abc unifi-protect-backup -${VERBOSITY}
